# systemd/novlearn-backend.service
[Unit]
Description=Novlearn Backend API
After=network.target postgresql.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/opt/novlearn/backend
Environment="PATH=/opt/novlearn/backend/venv/bin"
ExecStart=/opt/novlearn/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8010
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

